# infrastructure/openvpn_statusgen/teamcheck.conf
respawn
 
start on runlevel [2345]
stop on runlevel [06]
 
script
su - teamcheck -c "/home/teamcheck/gen_status_loop.py 2>&1" >>  /var/log/teamcheck.log
end script
