# infrastructure/team_testimage_scripts/service.conf
respawn
 
start on runlevel [2345]
stop on runlevel [06]
 
script
su - service -c "/home/service/service.py 2>&1" >>  /var/log/myservice.log
end script
