# infrastructure/vpn/gamescripts/antidos/antiddos.conf
respawn
 
start on runlevel [2345]
stop on runlevel [06]
 
script
su - antiddos -s /bin/bash -c "/home/antiddos/antiddos_serv.py 2>&1" >>  /var/log/antiddos.log
end script
