# infrastructure/vpn/gamescripts/openclosenetwork/networkclosed.conf
respawn
 
start on runlevel [2345]
stop on runlevel [06]
 
script
su - networkclosed -s /bin/bash -c "/home/networkclosed/networkclosed_serv.py 2>&1" >>  /var/log/networkclosed.log
end script
